# tailscale-expiry-checker.timer
[Timer]
OnActiveSec=0
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
